# qfang/qfang_2shou.conf
[program:qfang_2shou]
command=python3 qfang_2shoufang.py
directory=/root/fangyuan/qfang
autostart=true
autorestart=true
environment=LANG="en_US.utf8", LC_ALL="en_US.UTF-8", LC_LANG="en_US.UTF-8"