# anjuke/anjuke_2shou.conf
[program:anjuke_2shou]
command=python3 anjuke_2shoufang.py
directory=/root/fangyuan/anjuke
autostart=true
autorestart=true
environment=LANG="en_US.utf8", LC_ALL="en_US.UTF-8", LC_LANG="en_US.UTF-8"