# lianjia/lianjia_2shou.conf
[program:lianjia_2shou]
command=python3 lianjia_2shoufang.py
directory=/root/fangyuan/lianjia
autostart=true
autorestart=true
environment=LANG="en_US.utf8", LC_ALL="en_US.UTF-8", LC_LANG="en_US.UTF-8"