# tongcheng/58_2shou.conf
[program:58_2shou]
command=python3 58_2shoufang.py
directory=/root/fangyuan/58
autostart=true
autorestart=true
environment=LANG="en_US.utf8", LC_ALL="en_US.UTF-8", LC_LANG="en_US.UTF-8"